# 98-joystick.rules
ACTION=="add", ATTRS{idVendor}=="0737", ATTRS{idProduct}=="4716", RUN+="/usr/bin/jscal-restore /dev/input/js0"
